# zsh/etc/cron.daily/zsh-scrub
#!/bin/bash

zpool scrub `zpool list | tr -s " " | cut -d " " -f1 | egrep -v "NAME"`
